# autostart/claude-kbd.desktop
[Desktop Entry]
Type=Application
Name=Claude OSK
Exec=sh -c 'python3 /home/deck/.local/bin/claude-kbd.py >/tmp/claude-kbd-out.log 2>&1'
X-GNOME-Autostart-enabled=true
