# autostart/handheld-kbd-swap.desktop
[Desktop Entry]
Type=Application
Name=Better Handheld Keyboard Swap
Exec=__BIN__/handheld-kbd-swap.sh
X-GNOME-Autostart-enabled=true
